# xiguapi/custom/target/linux/rockchip/armv8/base-files/etc/uci-defaults/99-zzxgp-wifi-autosetup
#!/bin/sh
# /etc/uci-defaults/99-zzxgp-wifi-autosetup
# 首次启动WiFi配置入口脚本
# 调用预设配置脚本，成功后自删除

LOGFILE="/tmp/uci-defaults-wifi.log"
echo "========================================" > "$LOGFILE"
echo "$(date): 开始执行预设WiFi配置" >> "$LOGFILE"

# 检查核心脚本是否存在且可执行
if [ ! -f "/usr/local/bin/setup_wifi_universal.sh" ]; then
    echo "错误：核心脚本 /usr/local/bin/setup_wifi_universal.sh 不存在。" >> "$LOGFILE"
    exit 1 # 返回非0，脚本将被保留以供检查
fi

if [ ! -x "/usr/local/bin/setup_wifi_universal.sh" ]; then
    echo "错误：核心脚本不可执行，请检查权限。" >> "$LOGFILE"
    exit 1
fi

# 执行核心配置脚本
echo "正在执行 /usr/local/bin/setup_wifi_universal.sh ..." >> "$LOGFILE"
if /usr/local/bin/setup_wifi_universal.sh; then
    echo "$(date): 预设WiFi配置成功完成。" >> "$LOGFILE"
    echo "详细执行日志请查看: /tmp/wifi_preset_setup.log" >> "$LOGFILE"
    
    # 成功，删除自身（遵循uci-defaults机制）
    rm -f /etc/uci-defaults/99-zzxgp-wifi-autosetup
    echo "此启动脚本已自动删除。" >> "$LOGFILE"
    exit 0
else
    CONFIG_EXIT_CODE=$?
    echo "$(date): 预设WiFi配置失败，退出码: $CONFIG_EXIT_CODE" >> "$LOGFILE"
    echo "请检查日志: /tmp/wifi_preset_setup.log" >> "$LOGFILE"
    echo "此脚本将被保留，下次启动时重试。" >> "$LOGFILE"
    exit $CONFIG_EXIT_CODE
fi
